# tools/grass-ci/crontab
*/5    *    *     *     *   nice /work/grass-ci/grass-ci.sh
